# deploy/systemd/ocserv-agent.service
[Unit]
Description=ocserv-agent - OpenConnect VPN Server Agent
Documentation=https://github.com/dantte-lp/ocserv-agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ocserv-agent
Group=ocserv-agent

# Binary
ExecStart=/etc/ocserv-agent/ocserv-agent --config /etc/ocserv-agent/config.yaml

# Restart policy
Restart=always
RestartSec=5s
StartLimitInterval=300s
StartLimitBurst=5

# Resource limits
MemoryLimit=256M
CPUQuota=50%
TasksMax=100

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/etc/ocserv-agent /var/log/ocserv-agent /var/backups/ocserv-agent
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true

# Capabilities (needed for occtl and systemctl)
AmbientCapabilities=CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_ADMIN

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ocserv-agent

# Graceful shutdown
TimeoutStopSec=30s
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
